# trunk/packages/multiwan/files/etc/hotplug.d/firewall/50-multiwan
#!/bin/sh -u

keepDir="/tmp/multiwan/ifaces"

local state="${ZONE}_${INTERFACE}"
local active
config_get active "$state" ifname

[ -z "$active" ] &&
  {
  local ipaddr gateway
  config_get ipaddr  "$INTERFACE" ipaddr
  config_get gateway "$INTERFACE" gateway

  mkdir -p "$keepDir" &&
  cat <<EOF > "$keepDir/../$INTERFACE.tmp" &&
device='$DEVICE'
address='$ipaddr'
gateway='$gateway'
EOF
  mv -f "$keepDir/../$INTERFACE.tmp" "$keepDir/$INTERFACE"
  echo "Interface '$INTERFACE' registered with dev '$DEVICE', ip '$ipaddr', gw '$gateway'"
  } 2>&1 | logger -t multiwan -p daemon.notice
