# branches/kamikaze/packages/madwimax/files/etc/init.d/madwimax
#!/bin/sh /etc/rc.common

START=40

wimaxIface="wimax0"

start() {
        test "`uci get network.wan.ifname`" == "$wimaxIface" || {
           uci set network.wan.ifname="$wimaxIface"
           uci set network.wan.proto=dhcp
           uci set network.wan.auto=0
           uci commit
           }
        
        /usr/sbin/madwimax -vvfd
}

stop() {
        /usr/bin/killall madwimax
        /bin/sleep 3
        return 0
}