# branches/kamikaze/packages/multiwan/files/etc/init.d/multiwan-routekeeper
#!/bin/sh /etc/rc.common

START=80

start()
  {
  /sbin/multiwan-routekeeper.sh &
  }

stop()
  {
  local pidfile='/var/run/multiwan-routekeeper.pid'
  test -f "$pidfile" &&
  kill "`cat "$pidfile"`"
  }