# packages/set-srodi/files/etc/hotplug.d/iface/set-srodi
#!/bin/sh -ux

interfaceName="$INTERFACE"

{
echo "`date` '$interfaceName' '$0' '$@'"
set
echo "##########################################"
} >> "/tmp/hotplug-srodi.$interfaceName.txt"
