# packages/set-srodi/files/etc/hotplug.d/iface/set-srodi
#!/bin/sh -u

test "$PROTO" == 'static' &&
/sbin/multiwan-routing.sh "$INTERFACE" "$ACTION" "`uci get network.$INTERFACE.ipaddr`" "`uci get network.$INTERFACE.gateway`" 2>&1 | logger -t multiwan -p daemon.notice
